# dev/systemd/mony-dev.target
[Unit]
Description=mony development stack
Wants=mony-db.service mony-backend.service mony-frontend.service

[Install]
WantedBy=multi-user.target
